# wpp-deployer-webhook.service
[Unit]
Description=WPP Deployer Webhook Server
Documentation=https://github.com/youruser/wpp-deployer
After=network.target
Wants=network.target

[Service]
Type=simple
User=wpp-deployer
Group=wpp-deployer
WorkingDirectory=/opt/wpp-deployer
ExecStart=/usr/local/bin/wpp-deployer listen --port 3000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wpp-deployer-webhook

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/wpp-deployer

# Environment variables (optional)
Environment=WPP_WEBHOOK_PORT=3000
Environment=WPP_WEBHOOK_SECRET=

[Install]
WantedBy=multi-user.target 